# duplicati/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables

# declare webserviceinterface, webserviceallowedhostname, serverdatafolder,

## Get the 'message' key from the user config options.
serverdatafolder=$(bashio::config 'serverdatafolder')

## Print the message the user supplied, defaults to "Hello World..."
# bashio::log.info "${message:="Hello World..."}"

bashio::log.info "${serverdatafolder:="serverdatafolder"}"

## Run your program
bashio::log.info "Service Start /usr/bin/run_duplicati"
exec /usr/bin/run_duplicati
